# systemd/system/archcanary.path
[Unit]
Description=Trigger AUR malware check after pacman transactions

[Path]
PathChanged=/var/log/pacman.log
Unit=aur-malware-check-onchange.service

[Install]
WantedBy=multi-user.target
