# systemd/system/archcanary.path
[Unit]
Description=Trigger archcanary scan after pacman transactions

[Path]
PathChanged=/var/log/pacman.log
Unit=archcanary-onchange.service

[Install]
WantedBy=multi-user.target
