# tests/fake_home/.config/autostart/evil.desktop
[Desktop Entry]
Type=Application
Name=EvilApp
Exec=/tmp/evil_backdoor --daemon
