# configs/audit-rules.conf
# Auditd rules — archcanary
# Installed to /etc/audit/rules.d/30-archcanary.rules by ./install.sh --system
# Edit via the archcanary GUI (Utilities → Edit audit rules) or directly as root.
# Apply changes: sudo systemctl restart auditd
#
# Rule format:
#   -w /path -p [r][w][x][a] -k keyname   (watch file/directory)
#   -a always,exit -F arch=b64 -S syscall -k keyname   (syscall audit)
# Search logs: sudo ausearch -k keyname | aureport -i

# ---------------------------------------------------------------------------
# Startup: delete prior rules and set buffer size
# ---------------------------------------------------------------------------
-D
-b 8192
-f 1

# ---------------------------------------------------------------------------
# Authentication and privilege escalation
# ---------------------------------------------------------------------------

# Privileged commands
-a always,exit -F path=/usr/bin/sudo   -F perm=x -F auid>=1000 -F auid!=-1 -k priv_cmd
-a always,exit -F path=/usr/bin/su     -F perm=x -F auid>=1000 -F auid!=-1 -k priv_cmd
-a always,exit -F path=/usr/bin/pkexec -F perm=x -F auid>=1000 -F auid!=-1 -k priv_cmd
-a always,exit -F path=/usr/bin/newgrp -F perm=x -F auid>=1000 -F auid!=-1 -k priv_cmd

# Login records
-w /var/log/faillog  -p wa -k logins
-w /var/log/lastlog  -p wa -k logins

# Identity files
-w /etc/passwd   -p wa -k identity
-w /etc/shadow   -p wa -k identity
-w /etc/group    -p wa -k identity
-w /etc/gshadow  -p wa -k identity

# Sudo configuration
-w /etc/sudoers    -p wa -k sudoers
-w /etc/sudoers.d/ -p wa -k sudoers

# ---------------------------------------------------------------------------
# Package management — Arch / AUR
# ---------------------------------------------------------------------------

# Pacman database writes (package install / remove / upgrade)
-w /var/lib/pacman/local -p wa -k pacman

# makepkg executions (AUR builds)
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/makepkg -k aur_build
-a always,exit -F arch=b32 -S execve -F path=/usr/bin/makepkg -k aur_build

# ---------------------------------------------------------------------------
# Kernel modules
# ---------------------------------------------------------------------------

-w /usr/bin/kmod    -p x -k module_load
-w /usr/sbin/insmod -p x -k module_load
-w /usr/sbin/rmmod  -p x -k module_load
-a always,exit -F arch=b64 -S init_module,finit_module,delete_module -k module_load

# ---------------------------------------------------------------------------
# System configuration changes
# ---------------------------------------------------------------------------

-w /etc/hosts      -p wa -k system_config
-w /etc/hostname   -p wa -k system_config
-w /etc/fstab      -p wa -k system_config
-w /etc/os-release -p wa -k system_config

# Cron (scheduled task tampering)
-w /etc/crontab      -p wa -k cron
-w /etc/cron.d/      -p wa -k cron
-w /etc/cron.daily/  -p wa -k cron
-w /etc/cron.weekly/ -p wa -k cron

# Systemd units (watch for new/modified services)
-w /etc/systemd/system/ -p wa -k systemd_units

# SSH config (monitor even if sshd is disabled — detects re-enablement)
-w /etc/ssh/sshd_config    -p wa -k sshd_config
-w /etc/ssh/sshd_config.d/ -p wa -k sshd_config

# ---------------------------------------------------------------------------
# Optional: lock rules after loading (requires reboot to modify)
# Uncomment only when rules are finalized.
# ---------------------------------------------------------------------------
# -e 2
