# packaging/archcanary.install
post_install() {
  # Seed Lynis's false-positive-suppression profile, but only if lynis is
  # actually installed and the target has no custom.prf of its own yet —
  # mirrors install.sh's runtime check, which package() cannot do at build time.
  if command -v lynis &>/dev/null && [[ ! -f /etc/lynis/custom.prf ]]; then
    install -Dm644 /usr/lib/archcanary/lynis-custom.prf /etc/lynis/custom.prf
    echo "  installed: /etc/lynis/custom.prf (Lynis false-positive suppressions — edit to enable)"
  fi

  echo ""
  echo "==> archcanary installed."
  echo ""
  echo "  Enable the automated system scan (runs as root — weekly + on boot + after pacman):"
  echo "    sudo systemctl enable --now archcanary.timer archcanary.path"
  echo ""
  echo "  Enable the user-scope scan and result notifier (run as your user):"
  echo "    systemctl --user enable --now archcanary-user.timer archcanary-notify.path"
  echo ""
  echo "  Run a manual full scan:"
  echo "    archcanary --refresh --full --all-time"
  echo ""
  echo "  This package installs the system-wide pieces (binaries, systemd units,"
  echo "  polkit policy, /etc allowlists). Per-user setup can't be done from here"
  echo "  since it has to run as YOUR user, not the pacman transaction — run:"
  echo "    archcanary --doctor"
  echo "  to see exactly what's still missing (yay/paru hooks, aurscan, lynis)"
  echo "  and the precise command to fix each one."
  echo ""
  echo "  Optional add-ons: aurscan (LLM PKGBUILD scanner), yay/paru (Lua/PreBuildCommand hooks)."
  echo "  See: https://github.com/musqz/archcanary"
  echo ""
}

post_upgrade() {
  post_install
}
